# etc/rc.local
#!/bin/sh
plymouth --quit
rfkill unblock bluetooth
