# cron/crontab
# At minute 0 of every hour
0 * * * * rsync ~/alby /secondary
0 * * * * rsync ~/git /secondary
